# etc/cron.daily/dnsmasqconfig
#!/bin/bash
# reactivate dnsfilter failsafe
/usr/local/bin/dnsmasqconfig
