# etc/network/if-up.d/dnsmasqconfig
#!/bin/sh -e
# or create a symlink here (easier)
/usr/local/bin/dnsmasqconfig
