# etc/network/if-up.d/dnsmasqconfig
#!/bin/sh -e
/usr/local/bin/dnsmasqconfig
